# dnsdist/rootfs/etc/services.d/dnsdist/run
#!/usr/bin/with-contenv bashio
CONFIG="/etc/dnsdist.conf"

# Run dnsmasq
bashio::log.info "Starting dnsdist..."
exec dnsdist -C "${CONFIG}" < /dev/null